# bot.service
[Unit]
Description=systemd service file for teapot <3

[Service]
WorkingPath=/var/bot
ExecStart=node index.js

[Install]
WantedBy=multi-user.target